# networklog.conf
[program:networklog]
command=python3 /home/pi/networklog/main.py
numprocs=1
redirect_stderr=true
stdout_logfile=/home/pi/logs/net.log
stderr_logfile=/home/pi/logs/net_err.log
user=root
