# networklog.conf
[program:maimaipy]
command=python3 /home/pi/networklog/main.py
numprocs=1
redirect_stderr=true
stdout_logfile=/home/pi/logs/net.log
user=root
